# package-list-backup.timer
[Unit]
Description=package-list-backup timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
